# scripts/upload-db.timer
[Unit]
Description=Daily upload of SQLite database to Discord webhook

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
